# mail/.config/systemd/user/offlineimap-oneshot.timer
[Unit]
Description=Background sync up for offlineimap

[Timer]
OnUnitActiveSec=1min 

[Install]
WantedBy=mail.target
